# packages/rpm/unitfiles/opendaylight.service
[Unit]
Description=OpenDaylight SDN Controller
Documentation=https://wiki.opendaylight.org/view/Main_Page http://www.opendaylight.org/
After=network.service

[Service]
Type=forking
ExecStart=/opt/opendaylight/bin/start
User=odl
Group=odl
SuccessExitStatus=143
LimitNOFILE=102400
Restart=on-failure

[Install]
WantedBy=multi-user.target
